# Spec_sense_sensor_files/dev/fd/3/run/systemd/transient/user-120.slice
# This is a transient unit file, created programmatically via the systemd API. Do not edit.
[Unit]
Description=User Slice of gdm
After=systemd-logind.service
After=systemd-user-sessions.service

[Slice]
TasksMax=10813
