# Spec_sense_sensor_files/dev/fd/3/run/systemd/generator/boot-efi.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2duuid-2635\x2d885A.service
After=systemd-fsck@dev-disk-by\x2duuid-2635\x2d885A.service

[Mount]
Where=/boot/efi
What=/dev/disk/by-uuid/2635-885A
Type=vfat
Options=umask=0077
