# Spec_sense_sensor_files/dev/fd/3/run/systemd/generator/local-fs.target.requires/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
Where=/
What=/dev/disk/by-uuid/80f4e5a2-a2b2-49d9-ae1b-9076f7013ba5
Type=ext4
Options=errors=remount-ro
